# crontab
# SPDX-FileCopyrightText: 2025 Luis Villa <luis@lu.is>
#
# SPDX-License-Identifier: BlueOak-1.0.0

# News Fixed scheduled PDF generation
# Runs every Saturday at 7 AM UTC (Friday 11 PM PT) for weekend buffer
# Format: minute hour day-of-month month day-of-week command

0 7 * * 6 /usr/local/bin/uv run python /app/code/src/scheduled_generate.py --notify 2>&1
